# supervisord.conf
; supervisor config file.
; Start => /usr/bin/python /usr/bin/supervisord -c /www/default/docs/supervisord.conf

[supervisord]
logfile=/go/src/TestTask/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user = root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:Api]
command                 = /go/bin/api
process_name            = %(program_name)s_%(process_num)02d
numprocs                = 1
autostart               = true
autorestart             = true
stdout_logfile          = /go/src/TestTask/logs/api.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = /go/src/TestTask/logs/api_error.log
stderr_logfile_maxbytes = 1MB
user                    = root


[program:Parser]
command                 = /go/bin/parser
process_name            = %(program_name)s_%(process_num)02d
numprocs                = 1
autostart               = true
autorestart             = true
stdout_logfile          = /go/src/TestTask/logs/parser.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = /go/src/TestTask/logs/parser_error.log
stderr_logfile_maxbytes = 1MB
user                    = root


[program:Agregator]
command                 = /go/bin/agregator
process_name            = %(program_name)s_%(process_num)02d
numprocs                = 1
autostart               = true
autorestart             = true
stdout_logfile          = /go/src/TestTask/logs/agregator.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = /go/src/TestTask/logs/agregator.log
stderr_logfile_maxbytes = 1MB
user                    = root
